# tags/etc/service/lsphp/run
#!/bin/sh

exec 2>&1
exec su-exec $WEBSERVER_USER_NAME /usr/local/lsws/lsphp74/bin/lsphp -b *:3000